# scripts/systemd/devgrab-scheduler.service
[Unit]
Description=DevGrabBot Scheduler - AI-воркер анализа заявок
After=network.target postgresql.service redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/BOTS/DevGrabBot
ExecStart=/usr/bin/python3 -m src.ai.worker
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
EnvironmentFile=/var/www/BOTS/DevGrabBot/.env

[Install]
WantedBy=multi-user.target
